# dom0/i3.d/i3lock.service
[Unit]
Description=ScreenLock on suspend
Before=sleep.target

[Service]
User=alj
Type=forking
Environment=Display=:0
ExecStart=